# lifestack.service
sudo tee /etc/systemd/system/lifestack.service > /dev/null << 'EOF'
[Unit]
Description=Lifestack - Unified Life Services API
After=network.target mongodb.service
Wants=mongodb.service

[Service]
Type=simple
User=james
Group=james
WorkingDirectory=/home/james/apps/lifestack
EnvironmentFile=/home/james/apps/lifestack/.env
ExecStart=/usr/bin/node index.js

Restart=always
RestartSec=10
StartLimitInterval=200
StartLimitBurst=5

StandardOutput=journal
StandardError=journal
SyslogIdentifier=lifestack

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/james/apps/lifestack

LimitNOFILE=65536
MemoryLimit=512M

[Install]
WantedBy=multi-user.target
EOF